# infra/systemd/aharadar-api.service.example
[Unit]
Description=Aha Radar API (User Service)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/aharadar
EnvironmentFile=%h/.config/aharadar/api.env
ExecStart=%h/.nvm/versions/node/v22.22.0/bin/node packages/api/dist/main.js
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
NoNewPrivileges=true

[Install]
WantedBy=default.target
